# debian/cpufreqd.service
[Unit]
Description=Daemon for automatic cpufreq control depending on system usage
Conflicts=power-profiles-daemon.service tuned.service tlp.service auto-cpufreq.service system76-power.service
After=local-fs.target

[Service]
Type=exec
ExecStart=/usr/bin/cpufreqd
KillMode=mixed
Restart=on-failure
OOMScoreAdjust=-100

[Install]
WantedBy=multi-user.target
